# Software/Raspberry PI/devel_isolated/rosbash/etc/catkin/profile.d/15.rosbash.bash
# generated from rosbash/env-hooks/15.rosbash.bash.em

. "/mnt/dav/Projekt/Software/Raspberry PI/src/ros/rosbash/rosbash"
