# Software/Raspberry PI/devel_isolated/rosbash/etc/catkin/profile.d/15.rosbash.zsh
# generated from rosbash/env-hooks/15.rosbash.zsh.em

. "/mnt/dav/Projekt/Software/Raspberry PI/src/ros/rosbash/roszsh"
